# deploy/systemd/dspy-agent.target
[Unit]
Description=DSPy Agent Infrastructure Stack
Documentation=https://github.com/dspy/orchestrator
After=docker.service
Requires=docker.service

[Install]
WantedBy=multi-user.target

